# supervisord.ini
[supervisord]
nodaemon=true

[program:uwsgi]
command=/usr/local/bin/uwsgi --ini /etc/uwsgi/uwsgi.ini
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

stopsignal=QUIT


[fcgi-program:uvicorn-stripe]
socket=unix:///tmp/app.sock
socket_mode=777
command=/<path-to-workdir>/venv/bin/uvicorn --fd 0 main:app
directory=/<path-to-workdir>/
numprocs=1
process_name=uvicorn-stripe-%(process_num)d
stderr_logfile=/var/log/uvicorn-stripe.err.log
stdout_logfile=/var/log/uvicorn-stripe.out.log
autostart=true
autorestart=true